# dot_config/autostart/private_launch-dunst.sh.desktop
[Desktop Entry]
Exec=/home/aegis/.config/dunst/launch.sh
Icon=dialog-scripts
Name=launch-dunst.sh
Path=
Type=Application
X-KDE-AutostartScript=true
