# dot_config/autostart/private_launch-polybar.sh.desktop
[Desktop Entry]
Exec=/home/aegis/.config/polybar/launch.sh
Icon=dialog-scripts
Name=launch-polybar.sh
Path=
Type=Application
X-KDE-AutostartScript=true
